# scripts/user-session-monitor.service
[Unit]
Description=User Session Monitor
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/user-session-monitor run -config /etc/user-session-monitor/config.yaml
WorkingDirectory=/etc/user-session-monitor
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 